# vsftpd/service/syslogd/run
#!/bin/sh

exec /usr/sbin/syslog-stdout
